# sonic-ai-bridge.service
[Unit]
Description=Sonic Pi AI Bridge
After=network.target sonic-pi-headless.service
Requires=sonic-pi-headless.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/music
ExecStart=/usr/bin/python3 -u /root/music/ai_bridge.py
Restart=on-failure
RestartSec=3
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin

[Install]
WantedBy=multi-user.target
